# systemd-units/flickr-uploader-go.timer
[Unit]
Description=Run flickr uploader go periodically

[Timer]
OnBootSec=5min
OnUnitInactiveSec=2h

[Install]
WantedBy=timers.target